# emqx/rootfs/etc/s6-overlay/s6-rc.d/init-emqx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: EMQX
# Configures EMQX
# ==============================================================================
mkdir -p \
  /data/emqx/data \
  /data/emqx/etc \
  /data/emqx/plugins

#!/usr/bin/with-contenv bash
set -e

echo "Patching EMQX config.yaml using yq..."

yq eval '
  .dashboard.auth.method = "basic" |
  .dashboard.auth.basic.username = "admin" |
  .dashboard.auth.basic.password = "secret123"
' -i /opt/emqx/etc/config.yaml

echo "EMQX config.yaml patched successfully."
